# scripts/com.voicetoclaude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicetoclaude</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/federico/Documents/ClaudeCode/voice-to-claude/.venv/bin/python</string>
        <string>-m</string>
        <string>src.app</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/federico/Documents/ClaudeCode/voice-to-claude</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/federico/Library/Logs/voice-to-claude/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/federico/Library/Logs/voice-to-claude/stderr.log</string>
</dict>
</plist>
